# launchd/com.derek.ceo-agent-service.reply-consumer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.derek.ceo-agent-service.reply-consumer</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd /Users/derek/Documents/Projects/ceo-agent-service/apps/local-service &amp;&amp; export PATH=/Users/derek/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:$PATH PYTHONPATH=. CEO_NOT_SEND_MESSAGE=0 CEO_LIVE_SEND_BLOCKERS_ACCEPTED=1 CEO_WORKER_DB=/Users/derek/Documents/Projects/ceo-agent-service/data/auto-reply.sqlite3 CEO_WORKSPACE=/Users/derek/Documents/memory CEO_CORPUS_DIR=/Users/derek/Documents/Projects/ceo-agent-service/corpus &amp;&amp; .venv/bin/python -m ceo_agent_service.cli consume --poll-interval-seconds 10 --db /Users/derek/Documents/Projects/ceo-agent-service/data/auto-reply.sqlite3 --workspace /Users/derek/Documents/memory --corpus-dir /Users/derek/Documents/Projects/ceo-agent-service/corpus</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/derek/Documents/Projects/ceo-agent-service</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/derek/Library/Logs/ceo-agent-service/reply-consumer.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/derek/Library/Logs/ceo-agent-service/reply-consumer.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/derek</string>
    <key>PATH</key>
    <string>/Users/derek/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>PYTHONPATH</key>
    <string>.</string>
    <key>CEO_WORKER_DB</key>
    <string>/Users/derek/Documents/Projects/ceo-agent-service/data/auto-reply.sqlite3</string>
    <key>CEO_WORKSPACE</key>
    <string>/Users/derek/Documents/memory</string>
    <key>CEO_CORPUS_DIR</key>
    <string>/Users/derek/Documents/Projects/ceo-agent-service/corpus</string>
    <key>CEO_NOT_SEND_MESSAGE</key>
    <string>0</string>
    <key>CEO_LIVE_SEND_BLOCKERS_ACCEPTED</key>
    <string>1</string>
  </dict>
</dict>
</plist>
