# launchd/com.ceo-agent-service.main.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ceo-agent-service.main</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string><![CDATA[service_root="${CEO_SERVICE_ROOT:-${HOME:?HOME must be set}/Documents/Projects/ceo-agent-service}"; cd "${service_root}" && export PATH="${HOME}/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:${PATH:-}" && export HOME="${CEO_SERVICE_HOME:-${HOME}}" && export PYTHONPATH="${PYTHONPATH:-.}" && export CODEX_HOME="${CODEX_HOME:-${HOME}/.codex}" && export CEO_CODEX_MODEL="${CEO_CODEX_MODEL:-gpt-5.5}" && export CEO_CODEX_MODEL_REASONING_EFFORT="${CEO_CODEX_MODEL_REASONING_EFFORT:-medium}" && export CEO_WORKSPACE="${CEO_WORKSPACE:-${HOME}/Documents/memory}" && export DWS_DISABLE_KEYCHAIN="${DWS_DISABLE_KEYCHAIN:-1}" && export DWS_KEYCHAIN_DIR="${DWS_KEYCHAIN_DIR:-${service_root}/data/dws-keychain}" && export CEO_DING_ROBOT_NAME="${CEO_DING_ROBOT_NAME:-磊哥}" && export CEO_WORKER_DB="${CEO_WORKER_DB:-${HOME}/Library/Application Support/ceo-agent-service/auto-reply.sqlite3}" && export CEO_CORPUS_DIR="${CEO_CORPUS_DIR:-${service_root}/data/corpus}" && export CEO_OKR_LIVE_SOURCE_COMMAND="${CEO_OKR_LIVE_SOURCE_COMMAND:-${service_root}/scripts/dingteam_okr_live_source.py --user-id {user_id} --period-label {period_label}}" && export CEO_MAX_BATCHES="${CEO_MAX_BATCHES:-1}" && export CEO_NOT_SEND_MESSAGE=0 && export CEO_LIVE_SEND_BLOCKERS_ACCEPTED=1 && exec .venv/bin/python -m app.cli service --host "${CEO_AUDIT_WEB_HOST:-127.0.0.1}" --port "${CEO_AUDIT_WEB_PORT:-8765}" --db "${CEO_WORKER_DB}" --workspace "${CEO_WORKSPACE}" --corpus-dir "${CEO_CORPUS_DIR}"]]></string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/ceo-agent-service-main.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/ceo-agent-service-main.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>PYTHONPATH</key>
    <string>.</string>
  </dict>
</dict>
</plist>
