# launchd/com.derek.ceo-agent-service.audit-web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.derek.ceo-agent-service.audit-web</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>cd /Users/derek/Documents/Projects/ceo-agent-service/apps/local-service &amp;&amp; export PYTHONPATH=. HOME=/Users/derek/Documents/memory CODEX_HOME=/Users/derek/.codex DWS_DISABLE_KEYCHAIN=1 DWS_KEYCHAIN_DIR="/Users/derek/Documents/memory/Library/Application Support/dws-cli" CEO_WORKER_DB=/Users/derek/Documents/Projects/ceo-agent-service/data/auto-reply.sqlite3 CEO_WORKSPACE=/Users/derek/Documents/memory CEO_CORPUS_DIR=/Users/derek/Documents/Projects/ceo-agent-service/corpus PATH=/Users/derek/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin &amp;&amp; exec .venv/bin/python -m ceo_agent_service.cli audit-web --host 127.0.0.1 --port 8765</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/derek</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/derek/Library/Logs/ceo-agent-service/audit-web.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/derek/Library/Logs/ceo-agent-service/audit-web.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/derek/Documents/memory</string>
    <key>CODEX_HOME</key>
    <string>/Users/derek/.codex</string>
    <key>PATH</key>
    <string>/Users/derek/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>DWS_DISABLE_KEYCHAIN</key>
    <string>1</string>
    <key>DWS_KEYCHAIN_DIR</key>
    <string>/Users/derek/Documents/memory/Library/Application Support/dws-cli</string>
    <key>PYTHONPATH</key>
    <string>.</string>
    <key>CEO_WORKER_DB</key>
    <string>/Users/derek/Documents/Projects/ceo-agent-service/data/auto-reply.sqlite3</string>
    <key>CEO_WORKSPACE</key>
    <string>/Users/derek/Documents/memory</string>
    <key>CEO_CORPUS_DIR</key>
    <string>/Users/derek/Documents/Projects/ceo-agent-service/corpus</string>
    <key>CEO_AUDIT_WEB_HOST</key>
    <string>127.0.0.1</string>
    <key>CEO_AUDIT_WEB_PORT</key>
    <string>8765</string>
    <key>CEO_DING_ROBOT_NAME</key>
    <string>磊哥</string>
  </dict>
</dict>
</plist>
